# systemd/picamctl.service
[Unit]
Description=picamctl - Raspberry Pi Camera Control
After=network.target
Wants=network.target

[Service]
Type=simple
User=aachten
WorkingDirectory=/home/aachten/picamctl
ExecStart=/usr/bin/python3 /home/aachten/picamctl/picamctl.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target

